# deploy/kinkeeper-web.service
# 이 파일은 deploy/install.sh 가 USER와 PROJECT_DIR을 자동으로 교체합니다.
# 직접 복사해서 쓸 경우 User=, WorkingDirectory=, ExecStart=, EnvironmentFile= 경로를 수정하세요.
[Unit]
Description=KinKeeper Admin Web
After=network-online.target docker.service
Wants=docker.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/uvicorn web.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
EnvironmentFile=__PROJECT_DIR__/.env

[Install]
WantedBy=multi-user.target
