# deploy/kinkeeper-bot.service
# 이 파일은 systemd --user 유닛입니다. deploy/install.sh 가 PROJECT_DIR을 자동으로 교체해
# ~/.config/systemd/user/ 에 설치합니다.
# 직접 복사해서 쓸 경우 WorkingDirectory=, ExecStart=, EnvironmentFile= 경로를 수정하세요.
[Unit]
Description=KinKeeper Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=5
EnvironmentFile=__PROJECT_DIR__/.env

[Install]
WantedBy=default.target
